# systemd/predmkt-scanner.service
[Unit]
Description=Prediction Market — Scanner Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=predmkt
Group=predmkt
WorkingDirectory=/opt/predmkt
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/predmkt/.venv/bin/python /opt/predmkt/scanner.py --loop --interval 300
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=predmkt-scanner

[Install]
WantedBy=multi-user.target
